# systemd/arma-attendance-api.service
[Unit]
Description=Arma Attendance API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arma-attendance
Group=arma-attendance
WorkingDirectory=/opt/arma-attendance
EnvironmentFile=/opt/arma-attendance/.env
ExecStart=/usr/bin/node apps/api/dist/index.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
